# kubernetes/v1.34.0/service/setup-hosts.service
[Unit]
Description=Ensure dynamic IP (eth1 preferred, fallback to eth0) is mapped to kubernetes-1-34-0 in /etc/hosts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  IP=$(ip -4 addr show eth1 2>/dev/null | awk "/inet / {print \$2}" | cut -d/ -f1 | head -n1); \
  if [ -z "$IP" ]; then \
    IP=$(ip -4 addr show eth0 2>/dev/null | awk "/inet / {print \$2}" | cut -d/ -f1 | head -n1); \
  fi; \
  if [ -z "$IP" ]; then \
    echo "Error: Failed to get IP from either eth1 or eth0"; \
    exit 1; \
  fi; \
  sed -i "/[[:space:]]kubernetes-1-34-0$/d" /etc/hosts; \
  echo "$IP kubernetes-1-34-0" >> /etc/hosts'
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target